# examples/public/launchd/com.market-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_ROOT__/scripts/start_dashboard.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/market-dashboard-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/market-dashboard-launchd.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MARKET_PROJECT_ROOT</key>
        <string>__PROJECT_ROOT__</string>
        <key>MARKET_VENV_PATH</key>
        <string>__VENV_PATH__</string>
        <key>MARKET_DASHBOARD_SSH_TARGET</key>
        <string>__SSH_TUNNEL_TARGET__</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>

