# examples/public/launchd/com.market-nightly-postcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-nightly-postcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_ROOT__/scripts/nightly_full_with_postcheck.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/market-nightly-postcheck.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/market-nightly-postcheck.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MARKET_PROJECT_ROOT</key>
        <string>__PROJECT_ROOT__</string>
        <key>MARKET_PYTHON_BIN</key>
        <string>__PYTHON_BIN__</string>
    </dict>
</dict>
</plist>

